# etc/init/memory-match-web-1.conf
## ----------------------------------------------------------------------------
#
# Bookmarks
#
# http://kvz.io/blog/2009/12/15/run-nodejs-as-a-service-on-ubuntu-karmic/
#
## ----------------------------------------------------------------------------

# info
description "memory-match.chilts.org - The old classic for Firefox OS"
author      "Andrew Chilton"

# respawn this task
start on runlevel [2345]
respawn
respawn limit 20 5
stop on shutdown

# allow opening of more than 1024 files
limit nofile 4096 4096

# set some environment variables
env NODE_ENV=production
env DOMAIN=memory-match.chilts.org
env SESSION_SECRET=mi8Sheiw

# the script itself
script

    # quit the script if something goes wrong
    set -e

    # run the webserver as ubuntu
    exec \
        sudo -E -u ubuntu \
        /home/ubuntu/.nvm/v0.8.18/bin/node \
        /home/ubuntu/src/chilts-fxos-memory-match/app.js 5000 2>&1 >> /var/log/memory-match/web-1.log

end script

## ----------------------------------------------------------------------------
